# deploy/systemd/quantbuild-quantlog-report.service
[Unit]
Description=QuantLog post-run (validate/summarize/score) for previous UTC day
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/quantbuild/quantbuildv1
EnvironmentFile=-/etc/quantbuild/quantbuild.env
ExecStart=/bin/bash /opt/quantbuild/quantbuildv1/scripts/vps/quantlog_nightly.sh
StandardOutput=append:/opt/quantbuild/quantbuildv1/logs/quantlog_nightly.log
StandardError=append:/opt/quantbuild/quantbuildv1/logs/quantlog_nightly.log

[Install]
WantedBy=multi-user.target
